# ops/systemd/adsb-streamnew-upstream-report.timer
[Unit]
Description=Run ADS-B Stream-New upstream report-only baseline every 15 minutes

[Timer]
OnBootSec=8min
OnUnitActiveSec=15min
AccuracySec=2min
Persistent=true
Unit=adsb-streamnew-upstream-report.service

[Install]
WantedBy=timers.target
