# ops/systemd/adsb-streamnew-recovery-orchestrator.service
[Unit]
Description=ADS-B Stream-New recovery orchestrator shadow decision
After=network-online.target adsb-streamnew-subsystems-status.service
Wants=network-online.target

[Service]
Type=oneshot
User=yuki
Group=yuki
WorkingDirectory=/home/yuki/projects/stream_v2
EnvironmentFile=-/home/yuki/projects/stream_v2/.state/env/adsb-streamnew.env
EnvironmentFile=-/home/yuki/projects/stream_v2/.state/env/adsb-streamnew-youtube-monitor.env
Environment=PYTHONPATH=/home/yuki/projects/stream_v2/src
ExecStart=/home/yuki/projects/stream_v2/bin/stream-new recovery-orchestrator
