# ops/systemd/adsb-streamnew-network-observer.service
[Unit]
Description=ADS-B Stream-New network observer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yuki
Group=yuki
WorkingDirectory=/home/yuki/projects/stream_v3
EnvironmentFile=-/etc/default/stream-v3-observability-monitor
Environment=PYTHONPATH=/home/yuki/projects/stream_v3/src
ExecStart=/usr/bin/python3 /home/yuki/projects/stream_v3/src/watchers/network_observer.py --loop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
