# ops/systemd/adsb-streamnew-youtube-monitor.service
[Unit]
Description=ADS-B Stream-New YouTube public/live health monitor
After=network-online.target adsb-streamnew-youtube-stream.service
Wants=network-online.target

[Service]
Type=oneshot
User=yuki
Group=yuki
EnvironmentFile=-/home/yuki/projects/stream_v2/.state/env/adsb-streamnew-youtube-monitor.env
Environment=WATCHDOG_SYSTEMCTL_USE_SUDO=1
Environment=PYTHONPATH=/home/yuki/projects/stream_v2/src
ExecStart=/home/yuki/projects/stream_v2/venv/bin/python3 /home/yuki/projects/stream_v2/src/watchers/youtube_watchdog.py
