# ops/systemd/adsb-streamnew-memory-status.service
[Unit]
Description=ADS-B Stream-New memory guardrail snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=yuki
Group=yuki
WorkingDirectory=/home/yuki/projects/stream_v3
EnvironmentFile=-/etc/default/stream-v3-observability-monitor
Environment=PYTHONPATH=/home/yuki/projects/stream_v3/src
ExecStart=/home/yuki/projects/stream_v3/bin/stream-prod memory-status
