# ops/systemd/adsb-streamnew-resource-memory.timer
[Unit]
Description=Run ADS-B Stream-New resource memory diagnostic every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=adsb-streamnew-resource-memory.service

[Install]
WantedBy=timers.target
