# ops/systemd/stream-v3-remote-recovery.timer
[Unit]
Description=Run stream_v3 remote k3s recovery check every minute

[Timer]
OnBootSec=2min
OnActiveSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Unit=stream-v3-remote-recovery.service

[Install]
WantedBy=timers.target
