# ops/systemd/stream-v3-wan-address-observer-burst.service
[Unit]
Description=stream_v3 high-cadence WAN address observer burst for the recurring morning TCP stall window
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=yuki
EnvironmentFile=-/etc/default/stream-v3-wan-observers
Environment=STREAM_V3_REPO_DIR=/opt/stream_v3
Environment=WAO_STATE_DIR=/opt/stream_v3/.state/wan-observer
Environment="WAO_ANCHORS=cloudflare_v4=1.1.1.1:443,google_v4=8.8.8.8:443,cloudflare_v6=[2606:4700:4700::1111]:443,google_v6=[2001:4860:4860::8888]:443"
Environment=WAO_ANCHOR_TIMEOUT_SEC=1.5
Environment=WAO_PUBLIC_IPV4_URL=https://api.ipify.org
ExecStart=/bin/bash -lc 'cd "$${STREAM_V3_REPO_DIR}" && exec /usr/bin/python3 "$${STREAM_V3_REPO_DIR}/ops/scripts/wan_address_observer.py" --sample-reason jst_0800_0820_burst --interval-sec 10 --duration-sec 1200'
