# ops/systemd/stream-v3-health-snapshot.timer
[Unit]
Description=stream_v3 observability health snapshot timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Unit=stream-v3-health-snapshot.service

[Install]
WantedBy=timers.target
