# ops/systemd/stream-v3-netlink-wan-event-observer.service
[Unit]
Description=stream_v3 netlink WAN route/address/link event observer
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=yuki
EnvironmentFile=-/etc/default/stream-v3-wan-observers
Environment=STREAM_V3_REPO_DIR=/opt/stream_v3
Environment=NETLINK_WAN_STATE_DIR=/opt/stream_v3/.state/wan-observer
ExecStart=/bin/bash -lc 'cd "$${STREAM_V3_REPO_DIR}" && exec /usr/bin/python3 "$${STREAM_V3_REPO_DIR}/ops/scripts/netlink_wan_event_observer.py"'
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
