# ops/systemd/stream-v3-persistent-anchor-observer.service
[Unit]
Description=stream_v3 persistent TCP anchor observer
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=yuki
EnvironmentFile=-/etc/default/stream-v3-wan-observers
Environment=STREAM_V3_REPO_DIR=/opt/stream_v3
Environment=WAO_STATE_DIR=/opt/stream_v3/.state/wan-observer
Environment=WAO_PERSISTENT_INTERVAL_SEC=15
Environment=WAO_PERSISTENT_TIMEOUT_SEC=2.0
Environment="WAO_PERSISTENT_ANCHORS=cloudflare_v4|1.1.1.1|443|cloudflare-dns.com|cloudflare-dns.com|AS13335,google_v4|8.8.8.8|443|dns.google|dns.google|AS15169,cloudflare_v6|2606:4700:4700::1111|443|cloudflare-dns.com|cloudflare-dns.com|AS13335,google_v6|2001:4860:4860::8888|443|dns.google|dns.google|AS15169"
Environment=WAO_PERSISTENT_TRIGGER_WAN_SNAPSHOT=1
Environment=WAO_PERSISTENT_WAN_SNAPSHOT_INTERVAL_SEC=5
Environment=WAO_PERSISTENT_WAN_SNAPSHOT_CYCLES=7
Environment=WAO_PERSISTENT_WAN_SNAPSHOT_COOLDOWN_SEC=120
Environment=WAO_PERSISTENT_TRIGGER_RTMPS_BURST=1
Environment=WAO_PERSISTENT_RTMPS_BURST_INTERVAL_SEC=5
Environment=WAO_PERSISTENT_RTMPS_BURST_DURATION_SEC=300
Environment=WAO_PERSISTENT_RTMPS_BURST_COOLDOWN_SEC=120
Environment=WAO_PUBLIC_IPV4_URL=https://api.ipify.org
Environment=WAO_ANCHOR_TIMEOUT_SEC=1.5
ExecStart=/bin/bash -lc 'cd "$${STREAM_V3_REPO_DIR}" && exec /usr/bin/python3 "$${STREAM_V3_REPO_DIR}/ops/scripts/persistent_tcp_anchor_observer.py"'
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
